# roles/result/files/result.service
---
[Unit]
Description=Result service
After=multi.user.target
Conflicts=getty@tty1.Service

[Service]
User=root
Type=simple
ExecStart=/usr/bin/result/views/app.js
StandardInput=tty-force

[Install]
WantedBy=multi-user.target
